# deploy/systemd/demande-materiel.service
[Unit]
Description=Demande Materiel Flask app (Gunicorn)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/demande-materiel
EnvironmentFile=/opt/demande-materiel/.env
ExecStart=/opt/demande-materiel/.venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8080 --timeout 120 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
